# systemd/autoedge-backup.service
[Unit]
Description=Encrypted off-host backup for AutoEdge Licensing Server
After=network-online.target autoedge-licensing.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/autoedge-backup.env
ExecStart=/opt/autoedge-licensing/scripts/backup_production.sh
Nice=10
IOSchedulingClass=idle
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/backups/autoedge-licensing /var/cache/autoedge-restic
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
