# systemd/autoedge-backup.timer
[Unit]
Description=Daily AutoEdge encrypted off-host backup

[Timer]
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=30m
Persistent=true
Unit=autoedge-backup.service

[Install]
WantedBy=timers.target
