# systemd/autoedge-licensing.service
[Unit]
Description=AutoEdge Licensing Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=autoedge
Group=autoedge
WorkingDirectory=/opt/autoedge-licensing
EnvironmentFile=/etc/autoedge-licensing.env
ExecStart=/opt/autoedge-licensing/.venv/bin/python -m autoedge_licensing.app
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/autoedge-licensing
ReadOnlyPaths=/etc/autoedge-licensing/keys
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
